# root/etc/systemd/system/avorionservers.target
[Unit]
Description=Avorion Server Targets
After=backup.mount
Wants=avorion@ds9server.service

[Install]
WantedBy=multi-user.target